# home/etc/sv/pa-applet/run
#!/bin/sh

sv start ~/graphical/pipewire-pulse || exit 1
exec pa-applet
