# home/etc/sv/polybar/run
#!/bin/sh

[ -z "${I3SOCK}" ] || exit 1

exec polybar primary 2>&1
