# home/etc/sv/polybar/run
#!/bin/sh

[ -S "${XDG_RUNTIME_DIR}"/i3/ipc-socket.* ] || exit 1

exec polybar primary 2>&1
