# home/etc/sv/blueman-applet/run
#!/bin/sh

pgrep bluetoothd > /dev/null || exit 1
exec blueman-applet
