# home/etc/sv/nm-applet/run
#!/bin/sh

pgrep NetworkManager > /dev/null || exit 1
exec nm-applet
