# home/etc/sv/pipewire-pulse/run
#!/bin/sh

sv start ~/graphical/pipewire || exit 1
exec pipewire-pulse
