# home/etc/sv/pipewire/run
#!/bin/sh
pgrep dbus > /dev/null || exit 1
exec pipewire
